# src/cron/crontab
* * * * * root /var/www/bin/console app:job:terminate-subscriptions
