# src/cron/crontab
0 1 * * * root /var/www/bin/console app:job:terminate-subscriptions
